# backend/supervisord/fetch.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord-fetch.pid

[unix_http_server]
file=/tmp/supervisor-fetch.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor-fetch.sock

# --- workers (Pattern B: limited self-heal) ---
# stopasgroup + killasgroup: terminate taskiq's multiprocessing children at process group level
# autorestart=unexpected + startretries=3: absorb transient failures, persistent bugs hit FATAL after 3 retries
# exitcodes=0: only exit 0 is "expected", everything else is unexpected

# dispatch worker はpipeline:dispatchのdispatch / sweep controlを担う。pool cap (5+5) に揃え
# 起動時 backlog の thundering herd で pool を即枯渇させないよう同時実行を 10 に固定する。
[program:dispatch]
command=taskiq worker --workers 1 --max-async-tasks 10 app.queue.brokers:broker_dispatch app.queue.tasks.acquisition app.queue.tasks.completion --ack-type when_executed
stopasgroup=true
killasgroup=true
autorestart=unexpected
startretries=3
exitcodes=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# collection workerはbroker_collectionの共有consumerとしてpipeline:acquisition / pipeline:completion
# の2 Streamを読み、DB / 外部HTTP混在task (acquire_source / scrape_html_body)を捌くため、
# 起動時 backlog の一斉実行 (thundering herd) で DB pool を飽和させないよう同時実行を 5 に制限する。
[program:collection]
command=taskiq worker --workers 1 --max-async-tasks 5 app.queue.brokers:broker_collection app.queue.tasks.acquisition app.queue.tasks.completion --ack-type when_executed
stopasgroup=true
killasgroup=true
autorestart=unexpected
startretries=3
exitcodes=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# eventlistener: capture PROCESS_STATE_FATAL only (NOT EXITED).
# Transient failures (recover within 1-3 retries) do not fire this listener.
# Only persistent bugs (3 retries elapsed -> FATAL) trigger SIGTERM to supervisord
# -> container exit -> docker restart loop becomes visible from outside.

[eventlistener:fail_fast]
command=python /app/supervisord/fail_fast.py
events=PROCESS_STATE_FATAL
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
