# backend/supervisord/fetch.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord-fetch.pid

[unix_http_server]
file=/tmp/supervisor-fetch.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor-fetch.sock

# --- workers (Pattern B: limited self-heal) ---
# stopasgroup + killasgroup: terminate taskiq's multiprocessing children at process group level
# autorestart=unexpected + startretries=3: absorb transient failures, persistent bugs hit FATAL after 3 retries
# exitcodes=0: only exit 0 is "expected", everything else is unexpected

[program:metadata]
command=taskiq worker app.brokers:broker_metadata app.collection.tasks --ack-type when_executed
stopasgroup=true
killasgroup=true
autorestart=unexpected
startretries=3
exitcodes=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:content]
command=taskiq worker app.brokers:broker_content app.collection.tasks --ack-type when_executed
stopasgroup=true
killasgroup=true
autorestart=unexpected
startretries=3
exitcodes=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# eventlistener: capture PROCESS_STATE_FATAL only (NOT EXITED).
# Transient failures (recover within 1-3 retries) do not fire this listener.
# Only persistent bugs (3 retries elapsed -> FATAL) trigger SIGTERM to supervisord
# -> container exit -> docker restart loop becomes visible from outside.

[eventlistener:fail_fast]
command=python /app/supervisord/fail_fast.py
events=PROCESS_STATE_FATAL
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
