# backend/supervisord/insights.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord-insights.pid

[unix_http_server]
file=/tmp/supervisor-insights.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor-insights.sock

# trend-discovery は日次 cron・fan-out なし・全集計を 1 session=最大 1 connection で実行。
# engine も pool 2/2=cap4 に縮小済 (WORKER_POOL_SIZING)。同時実行 2 は実並列 1 に対する余裕。
[program:trend-discovery]
command=taskiq worker --workers 1 --max-async-tasks 2 app.queue.brokers:broker_trend_discovery app.queue.tasks.trend_discovery --ack-type when_executed
stopasgroup=true
killasgroup=true
autorestart=unexpected
startretries=3
exitcodes=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# briefing は週次 cron だが category master (最大 ~12) ごとに subtask を fan-out する。
# LLM は rate_limit 律速・DB は short checkout。同時実行を pool cap (5+5) に揃える。
[program:briefing]
command=taskiq worker --workers 1 --max-async-tasks 10 app.queue.brokers:broker_briefing app.queue.tasks.briefing --ack-type when_executed
stopasgroup=true
killasgroup=true
autorestart=unexpected
startretries=3
exitcodes=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[eventlistener:fail_fast]
command=python /app/supervisord/fail_fast.py
events=PROCESS_STATE_FATAL
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
