# backend/supervisord/scheduler.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord-scheduler.pid

[unix_http_server]
file=/tmp/supervisor-scheduler.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor-scheduler.sock

# 5 つの cron scheduler (dispatch / trend_discovery / briefing / agent / maintenance) を
# 1 プロセスに統合する。各 scheduler は自分の broker へ kick するため routing は不変
# (app.queue.scheduler_entrypoint 参照)。full app import を 5→1 回に畳み RSS を削減する。
[program:scheduler]
command=python -m app.queue.scheduler_entrypoint
stopasgroup=true
killasgroup=true
autorestart=unexpected
startretries=3
exitcodes=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[eventlistener:fail_fast]
command=python /app/supervisord/fail_fast.py
events=PROCESS_STATE_FATAL
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
